# fanscript.service
[Unit]
Description=Fan control script

[Service]
ExecStart=/usr/bin/fanscript.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
